# conf/supervisor/apache.conf
[program:apache2]
command=/opt/docker/start-apache.sh
numprocs=1
autostart=true
autorestart=true