# conf/supervisor/mysql.conf
[program:mysqld]
command=/opt/docker/start-mysql.sh
numprocs=1
autostart=true
autorestart=true